# systemd/indexer.service
[Unit]
Description=Indexer daemon with Unix socket API
After=network.target indexer.socket
Requires=indexer.socket

[Service]
Type=simple
Environment="INDEXER_PATH=/" "INDEXER_NAME=root" "INDEXER_INCLUDE_HIDDEN=true" "INDEXER_SOCKET=/var/run/indexer.sock" "INDEXER_DB_PATH=/tmp/indexer.db" "INDEXER_INTERVAL=1h" "INDEXER_LISTEN_FLAG="
EnvironmentFile=-/etc/default/indexer
ExecStartPre=/bin/sh -c 'if [ ! -f /etc/default/indexer ]; then echo "indexer: /etc/default/indexer missing; using unit defaults"; fi'
ExecStart=/usr/local/bin/indexer --path ${INDEXER_PATH} --name ${INDEXER_NAME} --include-hidden=${INDEXER_INCLUDE_HIDDEN} --db-path ${INDEXER_DB_PATH} --socket-path ${INDEXER_SOCKET} --interval ${INDEXER_INTERVAL} ${INDEXER_LISTEN_FLAG}
Restart=on-failure
RuntimeDirectory=indexer
RuntimeDirectoryMode=0755
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp
ProtectHostname=true
ProtectClock=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
